# deploy/systemd/claw-openclaw.service
[Unit]
Description=Claw OpenClaw Shared Runtime
After=network.target

[Service]
Type=simple
User=claw-manager
Group=claw-manager
Environment=OPENCLAW_BIN=/usr/local/bin/openclaw
Environment=CLAW_RUNTIME_ROOT=/srv/claw/runtime
Environment=OPENCLAW_CONFIG_PATH=/srv/claw/runtime/openclaw/openclaw.json
ExecStart=/bin/sh -lc '"$OPENCLAW_BIN" gateway'
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=3
WorkingDirectory=/srv/claw/runtime/openclaw

[Install]
WantedBy=multi-user.target
