# deploy/systemd/claw-openclaw.service
[Unit]
Description=Claw OpenClaw Shared Runtime
After=network.target

[Service]
Type=simple
User=claw-manager
Group=claw-manager
EnvironmentFile=/etc/claw-workspace-manager.env
Environment=HOME=/opt/claw-workspace-manager/home
Environment=CLAW_RUNTIME_ROOT=/opt/claw-workspace-manager/home/claw/runtime
Environment=OPENCLAW_CONFIG_PATH=/opt/claw-workspace-manager/home/claw/runtime/openclaw/openclaw.json
ExecStart=/bin/sh -lc '"$OPENCLAW_BIN" gateway'
ExecReload=/bin/kill -USR1 $MAINPID
Restart=on-failure
RestartSec=3
WorkingDirectory=/opt/claw-workspace-manager/home/claw/runtime/openclaw

[Install]
WantedBy=multi-user.target
